# mdnsd.service
[Unit]
Description=DNS Service Discovery mDNSResponder

[Service]
Type=simple
ExecStart=/usr/sbin/mdnsd -debug

[Install]
WantedBy=multi-user.target
